# j-login.service
[Unit]
Description=J-Login
Requires=dbus.socket

[Service]
ExecStart=/usr/sbin/j-login

[Install]
WantedBy=graphical.target
